# com.github.copilot-wrapper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.github.copilot-wrapper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/torsten/github/github_copilot_openai_api_wrapper/.venv/bin/copilot-wrapper</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/torsten/github/github_copilot_openai_api_wrapper</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>Nice</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityBackgroundIO</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/torsten/Library/Logs/copilot-wrapper/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/torsten/Library/Logs/copilot-wrapper/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/torsten/github/github_copilot_openai_api_wrapper/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
